# config/rc.local
#!/bin/sh -e

systemctl start avahi-alias@ota.local

exit 0
